# crates/agent/systemd/vm-agent.service
[Unit]
Description=Easy VM Cloud Agent
Documentation=https://github.com/your-org/easy-vm-cloud
After=network.target libvirtd.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/easy-vm-cloud/agent
ExecStart=/opt/easy-vm-cloud/agent/bin/agent
Restart=on-failure
RestartSec=10s

# 环境变量
Environment="NODE_ID="
Environment="NODE_NAME="
Environment="SERVER_ADDRESS=http://localhost:3000"
Environment="HEARTBEAT_INTERVAL=30"
Environment="LOG_LEVEL=info"

# 安全设置
CapabilityBoundingSet=CAP_NET_ADMIN CAP_SYS_ADMIN
NoNewPrivileges=false
PrivateTmp=true

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vm-agent

[Install]
WantedBy=multi-user.target

